# deploy/ah.service
[Unit]
Description=agentic-hosting daemon
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=ah
Group=ah
SupplementaryGroups=docker
EnvironmentFile=-/etc/default/ah
ExecStart=/usr/local/bin/ah serve \
  --port 8080 \
  --db-path /var/lib/ah/ah.db \
  --master-key-path /var/lib/ah/master.key
Restart=always
RestartSec=5s
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ah

StateDirectory=ah

# Security hardening
UMask=0077
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
PrivateDevices=yes
RestrictNamespaces=yes
ProtectKernelLogs=yes
ReadWritePaths=/var/lib/ah /etc/traefik/dynamic/ah.yml /run/docker.sock /var/run/docker.sock

[Install]
WantedBy=multi-user.target
